# setup/etc/cron.d/www.example.org.cron
# m h  dom mon dow  user command
* * * * * root if [[ -f /var/www/www.example.org/cron.php ]]; then cd /var/www/www.example.org && php cron.php > /dev/null; fi
